# systemd/stremio-client.service
[Unit]
Description=Stremio Media Client
Wants=stremio-server.service network-online.target graphical.target
After=stremio-server.service network-online.target graphical.target

[Service]
Type=forking
Environment=DISPLAY=:0.0
Environment=XAUTHORITY=/home/__USER__/.Xauthority
ExecStart=/usr/bin/bash /home/__USER__/Stremio-RaspberryPi/src/client-launcher.sh


[Install]
WantedBy=default.target
